# seriallcd.service
[Unit]
Description=Serial-to-LCD Daemon
After=network.target local-fs.target

[Service]
Type=simple
User=seriallcd
Group=seriallcd

# Where your binary lives (example)
ExecStart=/usr/local/bin/seriallcd --run

# Give the service its own private RAM filesystem
TemporaryFileSystem=/run/serial_lcd_cache:size=100M,mode=0755

# Application may read/write ONLY this directory
ReadWritePaths=/run/serial_lcd_cache

# Lock everything else down hard
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
NoNewPrivileges=true

# Restart on crash
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
